# configs/luwu-hw-autoconf.service
[Unit]
Description=Luwu OS Hardware Auto-Config (CM4 new/old detection)
After=multi-user.target
Wants=network.target

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 /home/pi/luwu-os/configs/hardware_autoconf.py
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
