# configs/luwu-hw-autoconf.service
[Unit]
Description=Luwu OS Hardware Auto-Config (CM4 new/old detection)
After=local-fs.target luwu-splash.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 /opt/luwu-os/configs/hardware_autoconf.py
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
